# usr/lib/systemd/systemd-udevd.service.d/25-defaults-SLE.conf
# This drop-in contains specific defaults used by SUSE distributions.
#
# Do not edit this file, it will be overwritten on update.
#
# If you need to override one of the options, please create your own
# drop-in in /etc/systemd with a higher priority (preferably greater
# than 70).

[Service]
# Some stuff spawned by systemd-udevd need to access to the hostname
# (bsc#1178374)
ProtectHostname=no
