# usr/lib/systemd/system/user@0.service.d/20-defaults-SUSE.conf
[Service]
# Delegation for privileged services is supposed to enable all or
# listed resource controllers, which would incur performance penalties
# to *all* running services (bsc#954765 jsc#PM-2229).
Delegate=no
