# usr/lib/systemd/system/user@0.service.d/25-defaults-SLE.conf
[Service]
# Adding processes to multiple cgroup controllers causes performance
# regressions (bsc#954765)
Delegate=no
