# usr/lib/systemd/system/systemd-logind.service.d/20-defaults-SUSE.conf
# To keep backward compat with system using NIS, turn off the
# IP sandboxing.
[Service]
IPAddressDeny=
